# code-messenger/start/supervisor/config.ini
[program:messenger-consume]
command=php /Users/Ramatou/Documents/personalProj/symfonycast/code-messenger/start/bin/console messenger:consume async --time-limit=3600
user=Ramatou
numprocs=2
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
