# advancedPhpObject/docker/front/config/supervisord_logstream.conf
[program:logs]
command=/usr/bin/logstream
autostart=true
autorestart=true
stopsignal=HUP
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
